# ec2-production/canvas-course-generator.service
[Unit]
Description=Canvas Course Shell Generator (React)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/canvas-course-generator
ExecStart=/usr/bin/node start-production.cjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=canvas-course-generator
Environment=NODE_ENV=production
Environment=PORT=5000

[Install]
WantedBy=multi-user.target